# app.service
[Unit]
Description=Task2

[Service]
ExecStart=/home/vagrant/main
Type=simple

[Install]
WantedBy=multi-user.target
